# SharpInit/Units/unit-test-files/dependency-test-4.service
[Unit]
Description=dummy unit 4 set to fail

[Service]
ExecStart="C:\Users\Admin\source\repos\DummyUnit\DummyUnit\bin\Debug\DummyUnit.exf"