# SharpInit/Units/unit-test-files/notepad.service
[Unit]
Description=Notepad

[Service]
ExecStart=notepad.exe